# insights-post-upload.service
[Unit]
Description=Insights post-upload processing

[Service]
Type=oneshot
ExecStart=/usr/libexec/insights-post-upload
