# insights-post-upload.path
[Unit]
Description=Insights post-upload trigger

[Path]
PathChanged=/etc/insights-client/.lastupload
Unit=insights-post-upload.service
